# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true

[program:app]
command=python /download_liked_coubs.py
autostart=false
autorestart=false